# www-apps/UFL-WebAdmin-SiteDeploy-Web/files/ufl_webadmin_sitedeploy_web.initd
#!/sbin/runscript

UFL_WEBADMIN_SITEDEPLOY_WEB_USER=${UFL_WEBADMIN_SITEDEPLOY_WEB_USER:-apache}
UFL_WEBADMIN_SITEDEPLOY_WEB_GROUP=${UFL_WEBADMIN_SITEDEPLOY_WEB_GROUP:-apache}
UFL_WEBADMIN_SITEDEPLOY_WEB_SCRIPT=${UFL_WEBADMIN_SITEDEPLOY_WEB_SCRIPT:-/usr/bin/ufl_webadmin_sitedeploy_web_fastcgi.pl}
UFL_WEBADMIN_SITEDEPLOY_WEB_PIDFILE=${UFL_WEBADMIN_SITEDEPLOY_WEB_PIDFILE:-/var/run/fastcgi/${SVCNAME}.pid}
UFL_WEBADMIN_SITEDEPLOY_WEB_SOCKET=${UFL_WEBADMIN_SITEDEPLOY_WEB_SOCKET:-/var/run/fastcgi/${SVCNAME}.socket}
UFL_WEBADMIN_SITEDEPLOY_WEB_NPROC=${UFL_WEBADMIN_SITEDEPLOY_WEB_NPROC:-1}

depend() {
	need net
	use dns logger postgresql
	before apache2
}

start() {
	ebegin "Starting ${SVCNAME} (${UFL_WEBADMIN_SITEDEPLOY_WEB_SCRIPT})"
	exec start-stop-daemon --start \
		--chuid ${UFL_WEBADMIN_SITEDEPLOY_WEB_USER}:${UFL_WEBADMIN_SITEDEPLOY_WEB_GROUP} \
		--exec "${UFL_WEBADMIN_SITEDEPLOY_WEB_SCRIPT}" \
		-- -l "${UFL_WEBADMIN_SITEDEPLOY_WEB_SOCKET}" -n "${UFL_WEBADMIN_SITEDEPLOY_WEB_NPROC}" -p "${UFL_WEBADMIN_SITEDEPLOY_WEB_PIDFILE}" -d
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop \
		--pidfile "${UFL_WEBADMIN_SITEDEPLOY_WEB_PIDFILE}"
	eend $?
}
